# tiaadeals.service
[Unit]
Description=TiaaDeals Backend Service
Requires=tiaadeals.socket
After=network.target postgresql.service

[Service]
Type=simple
User=ubuntu
Group=www-data
WorkingDirectory=/var/www/tiaadeals-backend
Environment=NODE_ENV=production
Environment=PM2_HOME=/home/ubuntu/.pm2

# Create socket directory with correct permissions
ExecStartPre=/bin/mkdir -p /var/run
ExecStartPre=/bin/chown root:www-data /var/run
ExecStartPre=/bin/chmod 775 /var/run
ExecStartPre=/bin/rm -f /var/run/tiaadeals.sock

# Start the application
ExecStart=/usr/bin/pm2 start ecosystem.config.js --no-daemon

# Reload and stop commands
ExecReload=/usr/bin/pm2 reload ecosystem.config.js
ExecStop=/usr/bin/pm2 stop ecosystem.config.js

# Restart policy
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=append:/var/log/tiaadeals/service.log
StandardError=append:/var/log/tiaadeals/error.log

[Install]
WantedBy=multi-user.target 